# skills/memory-distiller/com.wade.memory-distiller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wade.memory-distiller</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/brandonpackard/.openclaw/workspace/skills/memory-distiller/memory_distiller.py</string>
        <string>--apply</string>
        <string>--limit</string>
        <string>20</string>
        <string>--hours</string>
        <string>72</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/brandonpackard/.openclaw/logs/memory-distiller.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/brandonpackard/.openclaw/logs/memory-distiller.log</string>
</dict>
</plist>
